# ansible/roles/skype/files/pulseaudio.service
[Unit]
Description=Pulseaudio

[Service]
ExecStart=/usr/bin/pulseaudio --daemonize=no
Restart=on-failure

[Install]
WantedBy=default.target
